# systemd/dvsa-checker.service
[Unit]
Description=DVSA Test Slot Checker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/dvsa-checker
ExecStart=/home/pi/dvsa-checker/venv/bin/python -m src.scheduler
Restart=always
RestartSec=30
 EnvironmentFile=/home/pi/dvsa-checker/.env

# Memory limit — kill and restart if Chrome leaks
 MemoryMax=800M

# Allow graceful shutdown
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
